# trash/app/conf.d/blog.conf
description "blog"
author "rev"


env NODE=/usr/bin/python
env LOG_FILE=/var/log/blog/stdout.log
env PID_FILE=/var/log/blog/blog.pid
#env USER=gus

start on (filesystem and net-device-up IFACE=eth0)
stop on shutdown


respawn
respawn limit 10 5


script
    chdir /data/blog/
    exec start-stop-daemon --start --make-pidfile --pidfile $PID_FILE --exec $NODE $SCRIPT >> $LOG_FILE 2>&1
end script


pre-start script
    touch $LOG_FILE
    chown root:$USER $LOG_FILE
    chmod 770 $LOG_FILE
    echo "[`date`] blog starting" >> $LOG_FILE
end script


pre-stop script
    rm -f $PID_FILE
    echo "[`date`] blog stopping" >> $LOG_FILE
end script

